# hotplug.d/iface/25-update6
#!/bin/sh

[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan6" ] && {
        logger -t interface-hotplug "force update6"
        /etc/init.d/update6 force
}

